# ansible/redis.upstart.j2
description     "Redis"

start on (started networking)

respawn
respawn limit 5 60

exec docker run --rm --net=host --name jaciredis redis
