# ansible/redis.upstart.j2
description     "Redis"

start on (started networking)

respawn
respawn limit 5 60

pre-stop exec docker rm -f jaciredis

exec docker run -v /var/log --rm --net=host --name jaciredis redis
